# ops/systemd/execution.service
[Unit]
Description=AVWAP Bot - Execution V2 (orders/positions)
After=network-online.target
Wants=network-online.target
StartLimitInterval=60
StartLimitBurst=10

[Service]
Type=simple
WorkingDirectory=/root/avwap_r3k_scanner

EnvironmentFile=/root/avwap_r3k_scanner/.env
Environment=TEST_MODE=0

ExecStart=/root/avwap_r3k_scanner/venv/bin/python -u /root/avwap_r3k_scanner/execution_v2/execution_main.py --poll-seconds 300

Restart=always
RestartSec=5

StandardOutput=journal
StandardError=journal

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false

[Install]
WantedBy=multi-user.target

