# ops/systemd/execution-restart.service
[Unit]
Description=AVWAP Bot - Restart Execution After Daily Scan

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart execution.service
