# deploy/systemd/execution.service.d/10-runtime-dir.conf
[Service]
# Ensure /run/avwap exists for the flock lockfile (created each boot)
RuntimeDirectory=avwap
RuntimeDirectoryMode=0755
