# ops/systemd/execution.service.d/10-watchlist-gate.conf
[Service]
# Do not start execution unless today's watchlist exists (NY date)
ExecStartPre=/root/avwap_r3k_scanner/bin/check_watchlist_today.sh
