# ops/systemd/execution.service.d/20-restart-policy.conf
[Service]
# Avoid restart storms when the watchlist gate fails:
# - gate exits 1 (missing) or 2 (not today's NY date)
Restart=on-failure
RestartPreventExitStatus=1 2
